# doc/config-examples/systemd/olaris.service
[Unit]
Description=Olaris Media Server

[Service]
ExecStart=/opt/olaris/bin/olaris

[Install]
WantedBy=default.target
